# services/website.timer
[Timer]
OnBootSec=10 m
OnCalendar=daily
Unit=website.service

[Install]
WantedBy=timers.target